# systemd/edgecontroller-firstboot.service
[Unit]
Description=EdgeController first-boot hostname setup
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/edgecontroller-firstboot.sh
RemainAfterExit=yes
TimeoutStartSec=30

[Install]
WantedBy=multi-user.target
